# configs/dot-config/autostart/fcitx5.desktop
[Desktop Entry]
Type=Application
Name=Fcitx 5
Exec=fcitx5
Terminal=false
X-GNOME-Autostart-enabled=true
